# templates/default/confluence.service.erb
# Inspired by https://github.com/bparsons/atlassian-systemd/blob/master/confluence.service

[Unit] 
Description=Confluence Team Collaboration Software
After=network.target

[Service] 
Type=forking
PIDFile=<%= node[:confluence][:pid_dir] %>/catalina.pid
ExecStart=<%= node[:confluence][:install_path] %>/current/bin/startup.sh
ExecStop=<%= node[:confluence][:install_path] %>/current/bin/shutdown.sh
User=<%= node[:confluence][:run_user] %>

[Install] 
WantedBy=multi-user.target 
